# get-weekly-wx/get-weekly-wx.timer
[Unit]
Description=Run WEATHERAmerica weather forecast generator every Sunday at 7:00 AM

[Timer]
OnCalendar=Sun *-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
